# .codetmp/ecshopx-api/docker-new/supervisor/super-websocket.ini
[program:websocket]
command=/data/httpd/artisan websocket:start
stdout_logfile=/data/httpd/storage/logs/supervisor-websocket.log
redirect_stderr=true
autostart=true
autorestart=true
numprocs=1
user=root
startsecs=3
