# ecshopx-api/docker-new/supervisor/super-queue.ini
[program:bloated-queue-default]
command=/data/httpd/artisan doctrine:queue:work --queue=default --delay=3 --memory=128 --timeout=30 --sleep=1 --tries=3
stdout_logfile=/data/httpd/storage/logs/supervisor-bloated-queue-default.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=3
user=root
startsecs=3
startretries=100000

[program:bloated-queue-qucik]
command=/data/httpd/artisan doctrine:queue:work --queue=qucik --delay=3 --memory=128 --timeout=30 --sleep=1 --tries=3
stdout_logfile=/data/httpd/storage/logs/supervisor-bloated-queue-qucik.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=5
user=root
startsecs=3
startretries=100000

[program:bloated-queue-seckill]
command=/data/httpd/artisan doctrine:queue:work --queue=seckill --delay=3 --memory=128 --timeout=30 --sleep=1 --tries=3
stdout_logfile=/data/httpd/storage/logs/supervisor-bloated-queue-seckill.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=2
user=root
startsecs=3
startretries=100000

[program:bloated-queue-slow]
command=/data/httpd/artisan doctrine:queue:work --queue=slow --delay=3 --memory=128 --timeout=1800 --sleep=1 --tries=3
stdout_logfile=/data/httpd/storage/logs/supervisor-bloated-queue-slow.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=4
user=root
startsecs=3
startretries=100000

[program:bloated-queue-sms]
command=/data/httpd/artisan doctrine:queue:work --queue=sms --delay=3 --memory=128 --timeout=1800 --sleep=1 --tries=3
stdout_logfile=/data/httpd/storage/logs/supervisor-bloated-queue-sms.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=2
user=root
startsecs=3
startretries=100000
